# systemd-services/random-wallpaper/random-wallpaper.service
[Unit]
Description=Set a random wallpaper using swww

[Service]
ExecStart=%h/dotfiles/scripts/random-wallpaper.sh